# docker/supervisor-imghost.conf
[program:imghost]
command=/usr/local/bin/imghost
autorestart=true
directory=/usr/share/imghost
user=www-data
